# data/rpmmultiservice.service.in
[Unit]
Description=rpmmulti service

[Service]
Type=dbus
BusName=org.freedesktop.RpmMultiService
ExecStart=@rpmmultiserviceprivdir@/rpmmultiserviced
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target